# systemd/research-stream-sync.service
[Unit]
Description=Research Stream: pull Scholar alerts, enrich, fetch open-access PDFs
Documentation=https://github.com/snunezcr/paper-grabber
# Without this the job runs at boot before the network is up and every
# lookup fails.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# The IMAP app password lives here, not in the unit: unit files are
# world-readable, and `systemctl cat` would print the password.
# Create it with:  install -m 600 /dev/null %h/.config/paper-grabber/env
#   PAPER_GRABBER_IMAP_USER=you@gmail.com
#   PAPER_GRABBER_IMAP_PASSWORD=xxxx xxxx xxxx xxxx
#   OPENALEX_MAILTO=you@gmail.com
EnvironmentFile=%h/.config/paper-grabber/env

ExecStart=%h/Development/paper-grabber/scripts/daily.sh

# A laptop can suspend mid-run; let the job finish or fail rather than hang
# holding the ledger.
TimeoutStartSec=30min

# The service needs only its own data, the repo, and the network.
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/paper-grabber %h/.cache/paper-grabber
NoNewPrivileges=true

[Install]
WantedBy=default.target
